# archlinux/violante-server.service
[Unit]
Description=Violante Server

[Service]
ExecStart=/usr/bin/violante-server
Restart=always

[Install]
WantedBy=multi-user.target